# chef/cookbooks/custom_stuff/files/default/rc.local
#!/bin/sh

python /usr/local/bin/code.py >/tmp/code.log 2>&1 &

exit 0
